# vigilia-m4-gds3725/launchd/com.vigilia.matia-gds.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vigilia.matia-gds</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/vigilia/VIGILIA-MatIA/vigilia-m4-gds3725/scripts/matia_gds_service.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/vigilia/VIGILIA-MatIA/vigilia-m4-gds3725</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/vigilia/VIGILIA-MatIA/vigilia-m4-gds3725/runtime/logs/launchd-out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/vigilia/VIGILIA-MatIA/vigilia-m4-gds3725/runtime/logs/launchd-err.log</string>
</dict>
</plist>
